# indexer/supervisor.conf
[program:indexer]
command=poetry run indexer
directory=/app
user=indexer
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stopsignal=TERM

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid 